# blinkenlights.conf
[program:blinkenlights]
command=python /home/pi/code/blinkenlights/teststrip.py
user=pi
